# example/django/config/supervisor/uwsgi.conf
[program:uwsgi]
command = /usr/local/bin/uwsgi --ini /etc/uwsgi/pye-web-uwsgi.ini

user=ubuntu
group=ubuntu
